# forcetv/ftds/container-files/data/conf/supervisor.d/ftds.conf
[program:ftds]
command=/data/workbench/%(program_name)s_3%(process_num)s/ftds3.8.1.3-2018-01-25
process_name=%(program_name)s_3%(process_num)s
user=root
numprocs=PROCESS_COUNT
numprocs_start=PROCESS_START_NUMBER
directory=/data/workbench/%(program_name)s_3%(process_num)s/
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/data/logs/%(program_name)s_3%(process_num)s.log